# src/dashboard-data-backup/src/crontab
# Run backup script at 2am, 8am, 2pm, and 8pm UTC
0 2,8,14,20 * * * root cd /app/src && python backup_kusto_query_res.py >> /var/log/backup/backup.log 2>&1
# Empty line at end is required for cron
